# Library/LaunchAgents/ram_disk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>Create RAM disk</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/user/.local/bin/ramdisk</string>
    <string>512</string>
  </array>
  <key>LaunchOnlyOnce</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
